# x11/.config/autostart/clipit.desktop
[Desktop Entry]
Exec=clipit
Name=clipit
Comment=Lightweight GTK+ Clipboard Manager
Type=Application
