# x11/.config/autostart/emacs.desktop
[Desktop Entry]
Exec=/usr/bin/emacs --daemon
Name=emacs
Comment=GNU project Emacs
Type=Application
NotShowIn=EXWM
