# x11/.config/autostart/emacs.desktop
[Desktop Entry]
Exec=emacs --daemon
Name=emacs
Comment=GNU project Emacs
Type=Application
NotShowIn=EXWM
