# x11/.config/autostart/gnome-settings-daemon.desktop
[Desktop Entry]
Type=Application
Name=GNOME Settings Daemon
Exec=gsd-xsettings
NoDisplay=true
NotShowIn=KDE;GNOME;
