# x11/.config/autostart/onboard-autostart.desktop
[Desktop Entry]
Name=Onboard
Exec=onboard
Type=Application
NoDisplay=true
