# x11/.config/autostart/pasystray.desktop
[Desktop Entry]
Type=Application
Name=pasystray
Exec=pasystray
Comment=PulseAudio controller for the system tray
NoDisplay=true
