# x11/.config/autostart/pasystray.desktop
[Desktop Entry]
Exec=/usr/bin/pasystray
Name=pasystray
Comment=PulseAudio controller for the system tray
Type=Application
