# x11/.config/autostart/pcmanfm.desktop
[Desktop Entry]
Exec=/usr/bin/pcmanfm --daemon-mode
Name=pcmanfm
Comment=A lightweight Gtk+ based file manager for X Window.
Type=Application
