# x11/.config/autostart/pcmanfm.desktop
[Desktop Entry]
Exec=pcmanfm -d --desktop
Name=pcmanfm
Comment=A lightweight Gtk+ based file manager for X Window.
Type=Application
