# x11/.config/autostart/syncthing.desktop
[Desktop Entry]
Exec=/usr/bin/syncthing -no-browser -logflags=3
Name=syncthing
Comment=Syncthing
Type=Application
