# udev/etc/rules.d/100-android-mount.rules
ACTION=="add",SUBSYSTEM=="usb",ENV{ID_VENDOR_ID}=="0fce",ENV{ID_MODEL_ID}=="51a7",RUN+="/home/loic/bin/mtp-up"
