# x11/etc/systemd/system/xlock@.service
[Unit]
Description=Lock X session

[Service]
User=%I
Environment=DISPLAY=:0
ExecStart=/usr/bin/slock

[Install]
WantedBy=sleep.target
